# cron/dailyBuilds-hagfish.cron
6 10 * * * /bin/bash /Users/paul/github/dailyPackageBuilds/buildAll.sh > /dev/null

